# resources/cron.d/update
# update the mizmoz resources each night
0 0 * * * www-data php /path/to/director/mizmoz update